# systemd/opty-daemon.service
[Unit]
Description=opty — HDC-powered codebase context optimizer daemon
Documentation=https://github.com/boj/opty
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/opty global --port 7390
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
